# keboola-config/supervisord/services/app.conf
[program:app]
command=uv run streamlit run /src/app.py --server.port 8050 --server.headless true
directory=/src
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
